# var/lib/dpkg/info/awstats.conffiles
/etc/cron.d/awstats
/etc/default/awstats
/etc/awstats/awstats.conf
/etc/awstats/awstats.conf.local
