# var/lib/dpkg/info/fetchmail.list
/.
/usr
/usr/bin
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/doc
/usr/share/doc/fetchmail
/usr/share/doc/fetchmail/contrib
/usr/share/doc/fetchmail/examples
/usr/share/doc-base
/etc
/etc/ppp
/etc/ppp/ip-up.d
/etc/ppp/ip-down.d
/etc/resolvconf
/etc/resolvconf/update-libc.d
/etc/logcheck
/etc/logcheck/ignore.d.workstation
/etc/logcheck/ignore.d.server
/etc/default
/etc/init.d
/usr/bin/fetchmail
/usr/share/man/man1/fetchmail.1.gz
/usr/share/doc/fetchmail/contrib/fetchmail-mode.el.gz
/usr/share/doc/fetchmail/contrib/fetchsetup
/usr/share/doc/fetchmail/contrib/domino
/usr/share/doc/fetchmail/contrib/fetchmaildistrib
/usr/share/doc/fetchmail/contrib/mailqueue.pl.gz
/usr/share/doc/fetchmail/contrib/getmail
/usr/share/doc/fetchmail/contrib/gotmail
/usr/share/doc/fetchmail/contrib/gotmail.awk
/usr/share/doc/fetchmail/contrib/gotmail.conf
/usr/share/doc/fetchmail/contrib/gotmail.html.awk
/usr/share/doc/fetchmail/contrib/login
/usr/share/doc/fetchmail/contrib/logout
/usr/share/doc/fetchmail/contrib/maildaemon
/usr/share/doc/fetchmail/contrib/multidrop.gz
/usr/share/doc/fetchmail/contrib/README.gz
/usr/share/doc/fetchmail/contrib/novell
/usr/share/doc/fetchmail/contrib/poptest
/usr/share/doc/fetchmail/contrib/preauth-harness
/usr/share/doc/fetchmail/contrib/README.getmail
/usr/share/doc/fetchmail/contrib/toprocmail
/usr/share/doc/fetchmail/contrib/zsh-completion
/usr/share/doc/fetchmail/contrib/fetchmail.logrotate
/usr/share/doc/fetchmail/contrib/delete-later
/usr/share/doc/fetchmail/contrib/README.Debian
/usr/share/doc/fetchmail/contrib/fetchmailnochda.pl.gz
/usr/share/doc/fetchmail/contrib/runfetchmail.gz
/usr/share/doc/fetchmail/contrib/delete-later.README.gz
/usr/share/doc/fetchmail/README
/usr/share/doc/fetchmail/fetchmail-features.html
/usr/share/doc/fetchmail/design-notes.html
/usr/share/doc/fetchmail/todo.html
/usr/share/doc/fetchmail/fetchmail-FAQ.html
/usr/share/doc/fetchmail/OLDNEWS.gz
/usr/share/doc/fetchmail/changelog.Debian.gz
/usr/share/doc/fetchmail/copyright
/usr/share/doc/fetchmail/examples/fetchmailrc.example
/usr/share/doc/fetchmail/NEWS.Debian.gz
/usr/share/doc/fetchmail/NOTES.gz
/usr/share/doc/fetchmail/README.SSL.gz
/usr/share/doc/fetchmail/README.Debian.gz
/usr/share/doc-base/fetchmail
/etc/ppp/ip-up.d/fetchmail
/etc/ppp/ip-down.d/fetchmail
/etc/resolvconf/update-libc.d/fetchmail
/etc/logcheck/ignore.d.workstation/fetchmail
/etc/logcheck/ignore.d.server/fetchmail
/etc/default/fetchmail
/etc/init.d/fetchmail
/usr/bin/popclient
/usr/share/man/man1/fetchmailrc.1.gz
/usr/share/man/man1/popclient.1.gz
