# var/lib/dpkg/info/plesk-mail-pc-driver.conffiles
/etc/init.d/pc-remote
/etc/default/psa-pc-remote
/etc/rsyslog.d/plesk-mail.conf
