# var/lib/dpkg/info/plesk-mail-pc-driver.list
/.
/opt
/opt/psa
/opt/psa/handlers
/opt/psa/handlers/hooks
/opt/psa/admin
/opt/psa/admin/sbin
/opt/psa/admin/bin
/opt/psa/etc
/opt/psa/etc/modules
/opt/psa/etc/modules/watchdog
/opt/psa/etc/modules/watchdog/service.tpl.d
/var
/var/spool
/var/spool/postfix
/var/spool/postfix/plesk
/etc
/etc/init.d
/etc/default
/etc/rsyslog.d
/usr
/usr/lib
/usr/lib/plesk-9.0
/usr/lib/courier-imap
/usr/lib/courier-imap/authlib
/usr/lib/i386-linux-gnu
/usr/lib/i386-linux-gnu/sasl2
/usr/lib/postfix
/opt/psa/handlers/hooks/sender_ip
/opt/psa/handlers/hooks/check-quota
/opt/psa/admin/sbin/mailmng_domain_rename
/opt/psa/admin/sbin/listmng
/opt/psa/admin/sbin/mailmng
/opt/psa/admin/sbin/mailqueuemng
/opt/psa/etc/modules/watchdog/service.tpl.d/postfix
/etc/init.d/pc-remote
/etc/default/psa-pc-remote
/etc/rsyslog.d/plesk-mail.conf
/usr/lib/plesk-9.0/bouncemessage
/usr/lib/plesk-9.0/mailsrv_conf_init
/usr/lib/plesk-9.0/postfix-mailman
/usr/lib/plesk-9.0/mailsrv_entities_dump
/usr/lib/plesk-9.0/mailsrvanalog
/usr/lib/plesk-9.0/mailsrv_set_hostname
/usr/lib/plesk-9.0/postfix-poplockdb-clean
/usr/lib/plesk-9.0/postfix-queue
/usr/lib/plesk-9.0/mailmng_domain_toggle
/usr/lib/plesk-9.0/mail_admin_aliases
/usr/lib/plesk-9.0/mailman_lists_dump
/usr/lib/plesk-9.0/mail_postfix_transport_restore.worker
/usr/lib/plesk-9.0/postfix_conf_set
/usr/lib/plesk-9.0/psa-pc-remote
/usr/lib/plesk-9.0/postfix-local
/usr/lib/plesk-9.0/mail_postfix_transport_restore
/usr/lib/plesk-9.0/mailman_conf_init
/usr/lib/courier-imap/authlib/authpsa
/usr/lib/i386-linux-gnu/sasl2/libpleskauxprop.so
/usr/lib/postfix/plesk_saslauthd
/opt/psa/admin/bin/mailmng_domain_rename
/opt/psa/admin/bin/listmng
/opt/psa/admin/bin/mailmng
/opt/psa/admin/bin/mailqueuemng
