# var/lib/dpkg/info/plesk-base.list
/.
/opt
/opt/psa
/opt/psa/bin
/opt/psa/etc
/var
/var/lib
/var/lib/psa
/var/lib/psa/dumps
/var/lib/plesk
/var/lock
/var/lock/files
/etc
/etc/cron.daily
/etc/cron.weekly
/etc/psa
/etc/psa/private
/etc/sw
/etc/sw/keys
/etc/sw/keys/backup
/etc/sw/keys/keys
/etc/sw/keys/restart
/etc/sw/keys/lock
/etc/cron.monthly
/usr
/usr/lib
/usr/lib/plesk-9.0
/usr/lib/plesk-9.0/mysql
/opt/psa/bin/sw-engine-pleskrun
/opt/psa/etc/psa.conf.default
/etc/sw/keys/info
/etc/sw/keys/lock/ka
/etc/sw/keys/lock/repository
/usr/lib/plesk-9.0/mysql/mysql
