# var/lib/dpkg/info/psa-courier-imap.list
/.
/var
/var/run
/etc
/etc/courier-imap
/etc/courier-imap/shared.tmp
/etc/courier-imap/shared
/etc/init.d
/etc/pam.d
/usr
/usr/lib
/usr/lib/courier-imap
/usr/lib/courier-imap/md5
/usr/lib/courier-imap/numlib
/usr/lib/courier-imap/sha1
/usr/lib/courier-imap/libhmac
/usr/lib/courier-imap/authlib
/usr/sbin
/usr/bin
/usr/share
/usr/share/doc
/usr/share/doc/psa-courier-imap
/usr/share/courier-imap
/usr/share/man
/usr/share/man/man8
/usr/share/man/man7
/usr/share/man/man1
/etc/courier-imap/imapd.cnf
/etc/courier-imap/quotawarnmsg.example
/etc/courier-imap/pop3d
/etc/courier-imap/shared/index
/etc/courier-imap/pop3d.cnf
/etc/courier-imap/imapd
/etc/courier-imap/imapd-ssl
/etc/courier-imap/pop3d-ssl
/etc/init.d/courier-imap
/etc/pam.d/pop3
/etc/pam.d/imap
/usr/lib/courier-imap/md5/libmd5.a
/usr/lib/courier-imap/libauth.a
/usr/lib/courier-imap/numlib/libnumlib.a
/usr/lib/courier-imap/sha1/libsha1.a
/usr/lib/courier-imap/makedatprog
/usr/lib/courier-imap/pop3d-ssl.rc
/usr/lib/courier-imap/imapd-ssl.rc
/usr/lib/courier-imap/couriertcpd
/usr/lib/courier-imap/pop3d.rc
/usr/lib/courier-imap/libhmac/libhmac.a
/usr/lib/courier-imap/imapd.rc
/usr/lib/courier-imap/libauthmod.a
/usr/sbin/authenumerate
/usr/sbin/pop3login
/usr/sbin/sharedindexinstall
/usr/sbin/imaplogin
/usr/sbin/sharedindexsplit
/usr/sbin/userdbpw
/usr/sbin/courierlogger
/usr/bin/maildirmake
/usr/bin/maildirkw
/usr/bin/maildiracl
/usr/bin/pop3d
/usr/bin/deliverquota
/usr/bin/imapd
/usr/bin/couriertls
/usr/share/doc/psa-courier-imap/COPYING
/usr/share/doc/psa-courier-imap/INSTALL
/usr/share/doc/psa-courier-imap/README
/usr/share/doc/psa-courier-imap/COPYING.GPL
/usr/share/doc/psa-courier-imap/NEWS
/usr/share/doc/psa-courier-imap/ChangeLog
/usr/share/doc/psa-courier-imap/AUTHORS
/usr/share/courier-imap/userdb
/usr/share/courier-imap/mkimapdcert
/usr/share/courier-imap/makeuserdb
/usr/share/courier-imap/pw2userdb
/usr/share/courier-imap/vchkpw2userdb
/usr/share/courier-imap/mkpop3dcert
/usr/share/man/man8/mkimapdcert.8.gz
/usr/share/man/man8/userdbpw.8.gz
/usr/share/man/man8/makeuserdb.8.gz
/usr/share/man/man8/imapd.8.gz
/usr/share/man/man8/userdb.8.gz
/usr/share/man/man8/mkpop3dcert.8.gz
/usr/share/man/man8/deliverquota.8.gz
/usr/share/man/man7/authlib.7.gz
/usr/share/man/man1/maildiracl.1.gz
/usr/share/man/man1/courierlogger.1.gz
/usr/share/man/man1/couriertcpd.1.gz
/usr/share/man/man1/maildirmake.1.gz
/usr/share/man/man1/maildirkw.1.gz
/usr/sbin/userdb
/usr/sbin/mkimapdcert
/usr/sbin/makeuserdb
/usr/sbin/pw2userdb
/usr/sbin/vchkpw2userdb
/usr/sbin/mkpop3dcert
/usr/share/man/man8/pw2userdb.8.gz
/usr/share/man/man8/vchkpw2userdb.8.gz
/usr/share/man/man7/authdaemon.7.gz
/usr/share/man/man7/authmysql.7.gz
/usr/share/man/man7/authuserdb.7.gz
/usr/share/man/man7/authpam.7.gz
/usr/share/man/man7/authldap.7.gz
/usr/share/man/man7/authvchkpw.7.gz
/usr/share/man/man7/authdaemond.7.gz
/usr/share/man/man7/authcram.7.gz
/usr/share/man/man7/authpwd.7.gz
/usr/share/man/man7/authshadow.7.gz
