# var/lib/dpkg/info/psa-firewall.conffiles
/etc/init.d/psa-firewall-forward
/etc/init.d/psa-firewall
