# var/lib/dpkg/info/psa-vpn.list
/.
/opt
/opt/psa
/opt/psa/lib
/opt/psa/lib/modules
/opt/psa/lib/modules/vpn
/opt/psa/lib/modules/vpn/openvpn-windows
/opt/psa/lib/modules/vpn/openvpn-windows/system
/opt/psa/var
/opt/psa/var/modules
/opt/psa/var/modules/vpn
/opt/psa/admin
/opt/psa/admin/sbin
/opt/psa/admin/sbin/modules
/opt/psa/admin/sbin/modules/vpn
/opt/psa/admin/htdocs
/opt/psa/admin/htdocs/modules
/opt/psa/admin/htdocs/modules/vpn
/opt/psa/admin/htdocs/modules/vpn/images
/opt/psa/admin/bin
/opt/psa/admin/bin/modules
/opt/psa/admin/bin/modules/vpn
/opt/psa/admin/plib
/opt/psa/admin/plib/templates
/opt/psa/admin/plib/templates/modules
/opt/psa/admin/plib/templates/modules/vpn
/opt/psa/admin/plib/modules
/opt/psa/admin/plib/modules/vpn
/opt/psa/admin/plib/modules/vpn/ui
/opt/psa/admin/plib/modules/vpn/locales
/opt/psa/admin/plib/modules/vpn/locales/de-DE
/opt/psa/admin/plib/modules/vpn/locales/en-US
/opt/psa/admin/plib/modules/vpn/locales/zh-TW
/opt/psa/admin/plib/modules/vpn/locales/fr-FR
/opt/psa/admin/plib/modules/vpn/locales/ru-RU
/opt/psa/admin/plib/modules/vpn/locales/es-ES
/opt/psa/admin/plib/modules/vpn/locales/fi-FI
/opt/psa/admin/plib/modules/vpn/locales/ro-RO
/opt/psa/admin/plib/modules/vpn/locales/pt-BR
/opt/psa/admin/plib/modules/vpn/locales/ja-JP
/opt/psa/admin/plib/modules/vpn/locales/zh-CN
/opt/psa/admin/plib/modules/vpn/locales/pt-PT
/opt/psa/etc
/opt/psa/etc/modules
/opt/psa/etc/modules/vpn
/etc
/etc/init.d
/opt/psa/lib/modules/vpn/openvpn-windows/system/libssl32.dll
/opt/psa/lib/modules/vpn/openvpn-windows/system/libeay32.dll
/opt/psa/lib/modules/vpn/openvpn-windows/system/tapinstall.exe
/opt/psa/lib/modules/vpn/openvpn-windows/system/tap0801.sys
/opt/psa/lib/modules/vpn/openvpn-windows/system/OemWin2k.inf
/opt/psa/lib/modules/vpn/openvpn-windows/system/openvpn.exe
/opt/psa/lib/modules/vpn/openvpn-windows/Install TAP device.bat
/opt/psa/lib/modules/vpn/openvpn-windows/Uninstall TAP device.bat
/opt/psa/lib/modules/vpn/openvpn-windows/Connect to VPN.bat
/opt/psa/admin/sbin/modules/vpn/openvpn
/opt/psa/admin/sbin/modules/vpn/vpnctl
/opt/psa/admin/htdocs/modules/vpn/images/btn_generate_key_bg.gif
/opt/psa/admin/htdocs/modules/vpn/images/module-icon.gif
/opt/psa/admin/htdocs/modules/vpn/images/btn_download_pack_unix_bg.gif
/opt/psa/admin/htdocs/modules/vpn/images/btn_download_pack_win_bg.gif
/opt/psa/admin/htdocs/modules/vpn/index.php
/opt/psa/admin/htdocs/modules/vpn/module.css
/opt/psa/admin/bin/modules/vpn/zip
/opt/psa/admin/plib/templates/modules/vpn/vpn.tpl
/opt/psa/admin/plib/templates/modules/vpn/vpn_key_upload.tpl
/opt/psa/admin/plib/templates/modules/vpn/vpn_preferences.tpl
/opt/psa/admin/plib/modules/vpn/ui/modules.vpn.package.download.php
/opt/psa/admin/plib/modules/vpn/ui/modules.vpn.key.download.php
/opt/psa/admin/plib/modules/vpn/ui/modules.vpn.php
/opt/psa/admin/plib/modules/vpn/ui/modules.vpn.properties.php
/opt/psa/admin/plib/modules/vpn/ui/modules.vpn.package.php
/opt/psa/admin/plib/modules/vpn/ui/modules.vpn.key.php
/opt/psa/admin/plib/modules/vpn/ui/modules.vpn.key.upload.php
/opt/psa/admin/plib/modules/vpn/VPNMainForm.php
/opt/psa/admin/plib/modules/vpn/VPNConfigurator.php
/opt/psa/admin/plib/modules/vpn/locale.php
/opt/psa/admin/plib/modules/vpn/locales/de-DE/conhelp_de-DE.php
/opt/psa/admin/plib/modules/vpn/locales/de-DE/messages_de-DE.php
/opt/psa/admin/plib/modules/vpn/locales/en-US/conhelp_en-US.php
/opt/psa/admin/plib/modules/vpn/locales/en-US/messages_en-US.php
/opt/psa/admin/plib/modules/vpn/locales/zh-TW/conhelp_zh-TW.php
/opt/psa/admin/plib/modules/vpn/locales/zh-TW/messages_zh-TW.php
/opt/psa/admin/plib/modules/vpn/locales/fr-FR/messages_fr-FR.php
/opt/psa/admin/plib/modules/vpn/locales/fr-FR/conhelp_fr-FR.php
/opt/psa/admin/plib/modules/vpn/locales/ru-RU/messages_ru-RU.php
/opt/psa/admin/plib/modules/vpn/locales/ru-RU/conhelp_ru-RU.php
/opt/psa/admin/plib/modules/vpn/locales/es-ES/messages_es-ES.php
/opt/psa/admin/plib/modules/vpn/locales/es-ES/conhelp_es-ES.php
/opt/psa/admin/plib/modules/vpn/locales/fi-FI/messages_fi-FI.php
/opt/psa/admin/plib/modules/vpn/locales/fi-FI/conhelp_fi-FI.php
/opt/psa/admin/plib/modules/vpn/locales/ro-RO/conhelp_ro-RO.php
/opt/psa/admin/plib/modules/vpn/locales/ro-RO/messages_ro-RO.php
/opt/psa/admin/plib/modules/vpn/locales/pt-BR/conhelp_pt-BR.php
/opt/psa/admin/plib/modules/vpn/locales/pt-BR/messages_pt-BR.php
/opt/psa/admin/plib/modules/vpn/locales/ja-JP/messages_ja-JP.php
/opt/psa/admin/plib/modules/vpn/locales/ja-JP/conhelp_ja-JP.php
/opt/psa/admin/plib/modules/vpn/locales/zh-CN/conhelp_zh-CN.php
/opt/psa/admin/plib/modules/vpn/locales/zh-CN/messages_zh-CN.php
/opt/psa/admin/plib/modules/vpn/locales/pt-PT/messages_pt-PT.php
/opt/psa/admin/plib/modules/vpn/locales/pt-PT/conhelp_pt-PT.php
/opt/psa/admin/plib/modules/vpn/VPNUIPointer.php
/opt/psa/admin/plib/modules/vpn/VPNPreferencesForm.php
/opt/psa/admin/plib/modules/vpn/VPNKeyUploadForm.php
/opt/psa/etc/modules/vpn/vpn_db.sql
/etc/init.d/psa-vpn
/opt/psa/admin/bin/modules/vpn/openvpn
/opt/psa/admin/bin/modules/vpn/vpnctl
