# var/lib/dpkg/info/psa-firewall.list
/.
/opt
/opt/psa
/opt/psa/var
/opt/psa/var/modules
/opt/psa/var/modules/firewall
/opt/psa/admin
/opt/psa/admin/sbin
/opt/psa/admin/sbin/modules
/opt/psa/admin/sbin/modules/firewall
/opt/psa/admin/htdocs
/opt/psa/admin/htdocs/modules
/opt/psa/admin/htdocs/modules/firewall
/opt/psa/admin/htdocs/modules/firewall/images
/opt/psa/admin/bin
/opt/psa/admin/bin/modules
/opt/psa/admin/bin/modules/firewall
/opt/psa/admin/plib
/opt/psa/admin/plib/templates
/opt/psa/admin/plib/templates/modules
/opt/psa/admin/plib/templates/modules/firewall
/opt/psa/admin/plib/modules
/opt/psa/admin/plib/modules/firewall
/opt/psa/admin/plib/modules/firewall/ui
/opt/psa/admin/plib/modules/firewall/locales
/opt/psa/admin/plib/modules/firewall/locales/de-DE
/opt/psa/admin/plib/modules/firewall/locales/en-US
/opt/psa/admin/plib/modules/firewall/locales/zh-TW
/opt/psa/admin/plib/modules/firewall/locales/fr-FR
/opt/psa/admin/plib/modules/firewall/locales/ru-RU
/opt/psa/admin/plib/modules/firewall/locales/es-ES
/opt/psa/admin/plib/modules/firewall/locales/fi-FI
/opt/psa/admin/plib/modules/firewall/locales/ro-RO
/opt/psa/admin/plib/modules/firewall/locales/pt-BR
/opt/psa/admin/plib/modules/firewall/locales/ja-JP
/opt/psa/admin/plib/modules/firewall/locales/zh-CN
/opt/psa/admin/plib/modules/firewall/locales/pt-PT
/opt/psa/etc
/opt/psa/etc/modules
/opt/psa/etc/modules/firewall
/etc
/etc/init.d
/opt/psa/var/modules/firewall/firewall-emergency.sh
/opt/psa/admin/sbin/modules/firewall/safeact
/opt/psa/admin/htdocs/modules/firewall/images/move-down.gif
/opt/psa/admin/htdocs/modules/firewall/images/rule-partallow.gif
/opt/psa/admin/htdocs/modules/firewall/images/btn_revert_bg.gif
/opt/psa/admin/htdocs/modules/firewall/images/rule-deny.gif
/opt/psa/admin/htdocs/modules/firewall/images/module-icon.gif
/opt/psa/admin/htdocs/modules/firewall/images/rule-allow.gif
/opt/psa/admin/htdocs/modules/firewall/images/btn_edit_bg.gif
/opt/psa/admin/htdocs/modules/firewall/images/move-down-disabled.gif
/opt/psa/admin/htdocs/modules/firewall/images/btn_activate-disabled_bg.gif
/opt/psa/admin/htdocs/modules/firewall/images/btn_add_rule_bg.gif
/opt/psa/admin/htdocs/modules/firewall/images/move-up-disabled.gif
/opt/psa/admin/htdocs/modules/firewall/images/btn_revert-disabled_bg.gif
/opt/psa/admin/htdocs/modules/firewall/images/move-up.gif
/opt/psa/admin/htdocs/modules/firewall/images/btn_add_rule-disabled_bg.gif
/opt/psa/admin/htdocs/modules/firewall/images/btn_activate_bg.gif
/opt/psa/admin/htdocs/modules/firewall/index.php
/opt/psa/admin/htdocs/modules/firewall/module.css
/opt/psa/admin/plib/templates/modules/firewall/firewall.tpl
/opt/psa/admin/plib/templates/modules/firewall/rule_edit.tpl
/opt/psa/admin/plib/templates/modules/firewall/activate.tpl
/opt/psa/admin/plib/modules/firewall/ui/modules.firewall.edit.php
/opt/psa/admin/plib/modules/firewall/ui/modules.firewall.edit.rule.properties.php
/opt/psa/admin/plib/modules/firewall/ui/modules.firewall.edit.rule.toggle.php
/opt/psa/admin/plib/modules/firewall/ui/modules.firewall.edit.activate.php
/opt/psa/admin/plib/modules/firewall/ui/modules.firewall.edit.reset.php
/opt/psa/admin/plib/modules/firewall/ui/modules.firewall.edit.rule.php
/opt/psa/admin/plib/modules/firewall/ui/modules.firewall.php
/opt/psa/admin/plib/modules/firewall/ui/modules.firewall.edit.rule.delete.php
/opt/psa/admin/plib/modules/firewall/ui/modules.firewall.edit.swap-rules.php
/opt/psa/admin/plib/modules/firewall/RuleEditForm.php
/opt/psa/admin/plib/modules/firewall/locale.php
/opt/psa/admin/plib/modules/firewall/locales/de-DE/conhelp_de-DE.php
/opt/psa/admin/plib/modules/firewall/locales/de-DE/messages_de-DE.php
/opt/psa/admin/plib/modules/firewall/locales/en-US/conhelp_en-US.php
/opt/psa/admin/plib/modules/firewall/locales/en-US/messages_en-US.php
/opt/psa/admin/plib/modules/firewall/locales/zh-TW/conhelp_zh-TW.php
/opt/psa/admin/plib/modules/firewall/locales/zh-TW/messages_zh-TW.php
/opt/psa/admin/plib/modules/firewall/locales/fr-FR/messages_fr-FR.php
/opt/psa/admin/plib/modules/firewall/locales/fr-FR/conhelp_fr-FR.php
/opt/psa/admin/plib/modules/firewall/locales/ru-RU/messages_ru-RU.php
/opt/psa/admin/plib/modules/firewall/locales/ru-RU/conhelp_ru-RU.php
/opt/psa/admin/plib/modules/firewall/locales/es-ES/messages_es-ES.php
/opt/psa/admin/plib/modules/firewall/locales/es-ES/conhelp_es-ES.php
/opt/psa/admin/plib/modules/firewall/locales/fi-FI/messages_fi-FI.php
/opt/psa/admin/plib/modules/firewall/locales/fi-FI/conhelp_fi-FI.php
/opt/psa/admin/plib/modules/firewall/locales/ro-RO/conhelp_ro-RO.php
/opt/psa/admin/plib/modules/firewall/locales/ro-RO/messages_ro-RO.php
/opt/psa/admin/plib/modules/firewall/locales/pt-BR/conhelp_pt-BR.php
/opt/psa/admin/plib/modules/firewall/locales/pt-BR/messages_pt-BR.php
/opt/psa/admin/plib/modules/firewall/locales/ja-JP/messages_ja-JP.php
/opt/psa/admin/plib/modules/firewall/locales/ja-JP/conhelp_ja-JP.php
/opt/psa/admin/plib/modules/firewall/locales/zh-CN/conhelp_zh-CN.php
/opt/psa/admin/plib/modules/firewall/locales/zh-CN/messages_zh-CN.php
/opt/psa/admin/plib/modules/firewall/locales/pt-PT/messages_pt-PT.php
/opt/psa/admin/plib/modules/firewall/locales/pt-PT/conhelp_pt-PT.php
/opt/psa/admin/plib/modules/firewall/FwConfigurator.php
/opt/psa/admin/plib/modules/firewall/FirewallMainForm.php
/opt/psa/admin/plib/modules/firewall/RuleList.php
/opt/psa/admin/plib/modules/firewall/ActivateForm.php
/opt/psa/admin/plib/modules/firewall/BackendLinux.php
/opt/psa/admin/plib/modules/firewall/BackendFreeBSD.php
/opt/psa/admin/plib/modules/firewall/navigation.php
/opt/psa/admin/plib/modules/firewall/FirewallUIPointer.php
/opt/psa/etc/modules/firewall/firewall_db.sql
/etc/init.d/psa-firewall-forward
/etc/init.d/psa-firewall
/opt/psa/admin/bin/modules/firewall/safeact
