# var/lib/dpkg/info/quota.list
/.
/usr
/usr/share
/usr/share/man
/usr/share/man/man3
/usr/share/man/man1
/usr/share/man/man8
/usr/share/doc
/usr/share/doc/quota
/usr/share/quota
/usr/share/quota/ldap
/usr/share/lintian
/usr/share/lintian/overrides
/usr/bin
/usr/sbin
/var
/var/lib
/var/lib/quota
/sbin
/etc
/etc/init.d
/etc/cron.daily
/etc/default
/usr/share/man/man3/rquota.3.gz
/usr/share/man/man1/quota.1.gz
/usr/share/man/man8/rquotad.8.gz
/usr/share/man/man8/setquota.8.gz
/usr/share/man/man8/edquota.8.gz
/usr/share/man/man8/quotastats.8.gz
/usr/share/man/man8/warnquota.8.gz
/usr/share/man/man8/quotacheck.8.gz
/usr/share/man/man8/repquota.8.gz
/usr/share/man/man8/quota_nld.8.gz
/usr/share/man/man8/quotaon.8.gz
/usr/share/man/man8/quot.8.gz
/usr/share/man/man8/convertquota.8.gz
/usr/share/man/man8/xqmstats.8.gz
/usr/share/doc/quota/README.gettext
/usr/share/doc/quota/quotas.ms.gz
/usr/share/doc/quota/quotas.preformated.gz
/usr/share/doc/quota/README.Debian
/usr/share/doc/quota/quotas-1.eps.gz
/usr/share/doc/quota/README.ldap-support
/usr/share/doc/quota/warnquota.conf
/usr/share/doc/quota/README.mailserver
/usr/share/doc/quota/quotadoc.sgml.gz
/usr/share/doc/quota/changelog.Debian.gz
/usr/share/doc/quota/copyright
/usr/share/quota/ldap/applySystemQuotas.pl
/usr/share/quota/ldap/setquota-ldap.pl
/usr/share/quota/ldap/quota.schema
/usr/share/quota/ldap/setSystemQuotas.pl
/usr/share/quota/ldap/edquota_editor
/usr/share/lintian/overrides/quota
/usr/bin/quota
/usr/bin/quotasync
/usr/sbin/quot
/usr/sbin/setquota
/usr/sbin/quotastats
/usr/sbin/edquota
/usr/sbin/rpc.rquotad
/usr/sbin/convertquota
/usr/sbin/quota_nld
/usr/sbin/warnquota
/usr/sbin/xqmstats
/usr/sbin/repquota
/sbin/quotaon
/sbin/quotacheck
/etc/quotatab
/etc/warnquota.conf
/etc/init.d/quota
/etc/init.d/quotarpc
/etc/cron.daily/quota
/etc/default/quota
/etc/quotagrpadmins
/usr/share/man/man8/quotaoff.8.gz
/usr/share/man/man8/rpc.rquotad.8.gz
/sbin/quotaoff
