# var/lib/dpkg/info/sw-cp-server.conffiles
/etc/init.d/sw-cp-server
/etc/sw-cp-server/config
