# var/lib/dpkg/info/webalizer.conffiles
/etc/webalizer/webalizer.conf.sample
/etc/cron.daily/webalizer
