# var/lib/dpkg/info/webalizer.list
/.
/usr
/usr/bin
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/locale
/usr/share/locale/ca
/usr/share/locale/ca/LC_MESSAGES
/usr/share/locale/cs
/usr/share/locale/cs/LC_MESSAGES
/usr/share/locale/da
/usr/share/locale/da/LC_MESSAGES
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/el
/usr/share/locale/el/LC_MESSAGES
/usr/share/locale/es
/usr/share/locale/es/LC_MESSAGES
/usr/share/locale/et
/usr/share/locale/et/LC_MESSAGES
/usr/share/locale/fi
/usr/share/locale/fi/LC_MESSAGES
/usr/share/locale/fr
/usr/share/locale/fr/LC_MESSAGES
/usr/share/locale/gl
/usr/share/locale/gl/LC_MESSAGES
/usr/share/locale/hr
/usr/share/locale/hr/LC_MESSAGES
/usr/share/locale/hu
/usr/share/locale/hu/LC_MESSAGES
/usr/share/locale/id
/usr/share/locale/id/LC_MESSAGES
/usr/share/locale/is
/usr/share/locale/is/LC_MESSAGES
/usr/share/locale/it
/usr/share/locale/it/LC_MESSAGES
/usr/share/locale/ja
/usr/share/locale/ja/LC_MESSAGES
/usr/share/locale/ko
/usr/share/locale/ko/LC_MESSAGES
/usr/share/locale/lv
/usr/share/locale/lv/LC_MESSAGES
/usr/share/locale/ml
/usr/share/locale/ml/LC_MESSAGES
/usr/share/locale/nl
/usr/share/locale/nl/LC_MESSAGES
/usr/share/locale/no
/usr/share/locale/no/LC_MESSAGES
/usr/share/locale/pl
/usr/share/locale/pl/LC_MESSAGES
/usr/share/locale/pt
/usr/share/locale/pt/LC_MESSAGES
/usr/share/locale/ro
/usr/share/locale/ro/LC_MESSAGES
/usr/share/locale/ru
/usr/share/locale/ru/LC_MESSAGES
/usr/share/locale/sl
/usr/share/locale/sl/LC_MESSAGES
/usr/share/locale/sr
/usr/share/locale/sr/LC_MESSAGES
/usr/share/locale/sv
/usr/share/locale/sv/LC_MESSAGES
/usr/share/locale/tr
/usr/share/locale/tr/LC_MESSAGES
/usr/share/locale/uk
/usr/share/locale/uk/LC_MESSAGES
/usr/share/locale/zh
/usr/share/locale/zh/LC_MESSAGES
/usr/share/doc
/usr/share/doc/webalizer
/usr/share/doc/webalizer/examples
/etc
/etc/webalizer
/etc/cron.daily
/usr/bin/webalizer
/usr/share/man/man1/webalizer.1.gz
/usr/share/locale/ca/LC_MESSAGES/webalizer.mo
/usr/share/locale/cs/LC_MESSAGES/webalizer.mo
/usr/share/locale/da/LC_MESSAGES/webalizer.mo
/usr/share/locale/de/LC_MESSAGES/webalizer.mo
/usr/share/locale/el/LC_MESSAGES/webalizer.mo
/usr/share/locale/es/LC_MESSAGES/webalizer.mo
/usr/share/locale/et/LC_MESSAGES/webalizer.mo
/usr/share/locale/fi/LC_MESSAGES/webalizer.mo
/usr/share/locale/fr/LC_MESSAGES/webalizer.mo
/usr/share/locale/gl/LC_MESSAGES/webalizer.mo
/usr/share/locale/hr/LC_MESSAGES/webalizer.mo
/usr/share/locale/hu/LC_MESSAGES/webalizer.mo
/usr/share/locale/id/LC_MESSAGES/webalizer.mo
/usr/share/locale/is/LC_MESSAGES/webalizer.mo
/usr/share/locale/it/LC_MESSAGES/webalizer.mo
/usr/share/locale/ja/LC_MESSAGES/webalizer.mo
/usr/share/locale/ko/LC_MESSAGES/webalizer.mo
/usr/share/locale/lv/LC_MESSAGES/webalizer.mo
/usr/share/locale/ml/LC_MESSAGES/webalizer.mo
/usr/share/locale/nl/LC_MESSAGES/webalizer.mo
/usr/share/locale/no/LC_MESSAGES/webalizer.mo
/usr/share/locale/pl/LC_MESSAGES/webalizer.mo
/usr/share/locale/pt/LC_MESSAGES/webalizer.mo
/usr/share/locale/ro/LC_MESSAGES/webalizer.mo
/usr/share/locale/ru/LC_MESSAGES/webalizer.mo
/usr/share/locale/sl/LC_MESSAGES/webalizer.mo
/usr/share/locale/sr/LC_MESSAGES/webalizer.mo
/usr/share/locale/sv/LC_MESSAGES/webalizer.mo
/usr/share/locale/tr/LC_MESSAGES/webalizer.mo
/usr/share/locale/uk/LC_MESSAGES/webalizer.mo
/usr/share/locale/zh/LC_MESSAGES/webalizer.mo
/usr/share/doc/webalizer/changelog.gz
/usr/share/doc/webalizer/DNS.README.gz
/usr/share/doc/webalizer/README.gz
/usr/share/doc/webalizer/README.FIRST.gz
/usr/share/doc/webalizer/msfree.png
/usr/share/doc/webalizer/webalizer.png
/usr/share/doc/webalizer/README.Debian
/usr/share/doc/webalizer/copyright
/usr/share/doc/webalizer/examples/sample.conf.gz
/usr/share/doc/webalizer/examples/track_hist.gz
/usr/share/doc/webalizer/changelog.Debian.gz
/usr/share/doc/webalizer/country-codes.txt.gz
/etc/webalizer/webalizer.conf.sample
/etc/cron.daily/webalizer
/usr/bin/webazolver
/usr/share/man/man1/webazolver.1.gz
