# var/lib/dpkg/info/wide-dhcpv6-client.list
/.
/var
/var/lib
/var/lib/dhcpv6
/etc
/etc/wide-dhcpv6
/etc/init.d
/etc/network
/etc/network/if-down.d
/etc/network/if-up.d
/usr
/usr/sbin
/usr/share
/usr/share/doc
/usr/share/doc/wide-dhcpv6-client
/usr/share/doc/wide-dhcpv6-client/examples
/usr/share/man
/usr/share/man/man8
/usr/share/man/man5
/etc/wide-dhcpv6/dhcp6c.conf
/etc/wide-dhcpv6/dhcp6c-script
/etc/wide-dhcpv6/dhcp6c-ifupdown
/etc/init.d/wide-dhcpv6-client
/usr/sbin/dhcp6c
/usr/sbin/dhcp6ctl
/usr/share/doc/wide-dhcpv6-client/copyright
/usr/share/doc/wide-dhcpv6-client/NEWS.Debian.gz
/usr/share/doc/wide-dhcpv6-client/examples/dhcp6c.conf.sample
/usr/share/doc/wide-dhcpv6-client/changelog.Debian.gz
/usr/share/man/man8/dhcp6c.8.gz
/usr/share/man/man8/dhcp6ctl.8.gz
/usr/share/man/man5/dhcp6c.conf.5.gz
/etc/network/if-down.d/wide-dhcpv6-client
/etc/network/if-up.d/wide-dhcpv6-client
