# var/lib/dpkg/info/xinetd.list
/.
/etc
/etc/xinetd.d
/etc/init
/etc/init.d
/etc/default
/usr
/usr/sbin
/usr/share
/usr/share/man
/usr/share/man/man5
/usr/share/man/man8
/usr/share/doc
/usr/share/doc/xinetd
/usr/share/doc/xinetd/examples
/etc/xinetd.d/daytime
/etc/xinetd.d/chargen
/etc/xinetd.d/echo
/etc/xinetd.d/discard
/etc/xinetd.d/time
/etc/init/xinetd.conf
/etc/default/xinetd
/etc/xinetd.conf
/usr/sbin/xinetd
/usr/sbin/itox
/usr/sbin/xconv.pl
/usr/share/man/man5/xinetd.log.5.gz
/usr/share/man/man5/xinetd.conf.5.gz
/usr/share/man/man8/xconv.pl.8.gz
/usr/share/man/man8/xinetd.8.gz
/usr/share/man/man8/itox.8.gz
/usr/share/doc/xinetd/NEWS.Debian.gz
/usr/share/doc/xinetd/README.gz
/usr/share/doc/xinetd/README.update-inetd
/usr/share/doc/xinetd/xinetd.org-FAQ.html
/usr/share/doc/xinetd/README.Debian
/usr/share/doc/xinetd/copyright
/usr/share/doc/xinetd/examples/sample.conf.gz
/usr/share/doc/xinetd/changelog.Debian.gz
/etc/init.d/xinetd
